# up/luci-app-udpxy/root/etc/hotplug.d/iface/99-udpxy
#!/bin/sh
if [ "$ACTION" = "ifup" -o  "$ACTION" = "ifupupdate"  ] ;then
  ip a | grep pppoe-  
/etc/init.d/udpxy restart
fi