# up/pass/luci-app-bypass/root/etc/hotplug.d/iface/21-bypass
#!/bin/sh
[ "$ACTION" = ifup ] || exit 0
if uci -q get bypass.@global[0].global_server >/dev/null;then
	/etc/init.d/bypass stop
	/etc/init.d/bypass start

fi
