# up/pass/luci-app-bypass/root/etc/hotplug.d/iface/21-bypass
#!/bin/sh
if [ "$ACTION" = ifup -a "$INTERFACE" != wan_6 ];then
	if uci -q get bypass.@global[0].global_server >/dev/null;then

		/etc/init.d/bypass restart
	fi
fi
