# up/luci-theme-xkucat/root/etc/uci-defaults/30_luci-theme-opentopd
#!/bin/sh

uci batch <<-EOF
	set luci.themes.kucat=/luci-static/kucat
	set luci.main.mediaurlbase=/luci-static/kucat
	commit luci
EOF
exit 0
