# up/luci-app-msd_lite/root/etc/hotplug.d/iface/99-msd_lite
#!/bin/sh
if [[ "$ACTION" == "ifup" || "$ACTION" == "ifupupdate" ]]; then
 if ip a | grep -q pppoe-; then
   /etc/init.d/msd_lite restart
 fi
fi